# systemd_units/client/sns_backup.service
[Unit]
Description=SNSBackup

[Service]
ExecStart=/usr/bin/python3 /home/sambryant/bin/github_projects/pi-server/src/client/sns_backup.py
